# playbooks/firefly/files/app.uploads-backup.path
[Unit]
Description=Checks if firefly uploads directory have changed

[Path]
PathChanged=/opt/firefly/uploads

[Install]
WantedBy=default.target