# roles/nfsv4-client-setup/files/nfs-kinit.timer
[Unit]
Description=Daily kinit of KDC ticket

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target