# playbooks/firefly/templates/app.uploads-backup.path
[Unit]
Description=Checks if firefly uploads directory have changed

[Path]
PathChanged={{ debian_additional_disks["appdata"]["mount_point"] }}/firefly/uploads

[Install]
WantedBy=default.target