# playbooks/firefly/templates/app.uploads-backup.path
[Unit]
Description=Checks if firefly uploads directory have changed

[Path]
PathChanged={{ docker_project_src }}/uploads

[Install]
WantedBy=default.target